# service/log/run
#!/bin/sh
exec 2>&1
exec multilog s1000000 n4 /var/log/dbus-hass-smartmeter